# venus-os/service/log/run
#!/bin/sh
exec multilog t s25000 n20 /var/log/victron-ble-dbus
